# os-builder/stage3-nomad/03-setup-services/files/nomad-pi.service
[Unit]
Description=Nomad Pi Server
After=network.target

[Service]
Type=simple
User=pi
# Must run from the repo root: app/main.py uses absolute `from app import ...`
# imports, which don't resolve when the working directory is app/ itself.
WorkingDirectory=/opt/nomad-pi
Environment="PATH=/opt/nomad-pi/venv/bin"
# Port 8000 matches setup.sh installs, the README, and mDNS URLs — and a
# non-root user cannot bind port 80 anyway.
ExecStart=/opt/nomad-pi/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --timeout-keep-alive 300
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
