# os-builder/stage3-nomad/03-setup-services/files/nomad-pi-firstboot.service
[Unit]
Description=Nomad Pi first-boot configuration (hotspot setup)
After=NetworkManager.service
Wants=NetworkManager.service
ConditionPathExists=!/etc/nomad-pi/.firstboot-done

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/nomad-pi-firstboot.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
